# deploy/macos/launchd/com.uw.watchtower.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.uw.watchtower.worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd __APP_DIR__ &amp;&amp; eval "$(mamba shell.bash hook)" &amp;&amp; mamba activate watchtower &amp;&amp; watchtower worker run --node-id __NODE_ID__</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__APP_DIR__</string>
    <key>StandardOutPath</key>
    <string>__DATA_ROOT__/logs/worker.log</string>
    <key>StandardErrorPath</key>
    <string>__DATA_ROOT__/logs/worker.err.log</string>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
